# archive/vpn-productd/systemd/vpn-product-lifecycle-daily.timer
[Unit]
Description=Run vpn-product lifecycle sync daily

[Timer]
OnCalendar=*-*-* 03:10:00
Persistent=true
Unit=vpn-product-lifecycle-daily.service

[Install]
WantedBy=timers.target
