# deploy/remnawave/systemd/vpn-backup.timer
[Unit]
Description=Daily Remnawave backup timer (00:00 UTC)

[Timer]
OnCalendar=*-*-* 00:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
